# admin/supervisor_snd.ini
[program:snd]
command=/home/www/.local/bin/gunicorn ndrive.wsgi
numprocs=1
directory=/home/www/Super-Neutron-Drive/server
stdout_logfile=/home/www/Super-Neutron-Drive/logs/server.log
redirect_stderr=true
autostart=true
autorestart=true
user=www
